# 001-aws-datadog/custom.service
[Unit]
Description=custom metrics job

[Service]
Type=oneshot
ExecStart=python3 /home/ubuntu/script.py
