# 001-aws-datadog/custom.timer
[Unit]
Description=custom metrics job timer

[Timer]
OnCalendar=*:1/1

[Install]
WantedBy=timers.target
